# scripts/debian/utilities/obsidian_sync/obsidian-sync.timer
[Unit]
Description=Periodic Obsidian vault sync
Wants=obsidian-sync.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
