# scripts/debian/utilities/obsidian_sync/obsidian-sync.service
[Unit]
Description=Obsidian vault inbound sync
After=network-online.target

[Service]
Type=oneshot
User=tbryant
EnvironmentFile=/home/tbryant/.hermes/.env
ExecStart=/usr/local/bin/obsidian-sync in
